# scripts/conf/swift_exporter.conf
start on startup

script
   /opt/ss/support/swift_exporter/bin/swift_exporter
end script
